# other/perfomance/perfomance.service
# /etc/systemd/system/perfomance.service

[Service]
ExecStart=/usr/local/bin/disk-space-check.sh